# etc/cron.daily/checkupdates
#!/usr/bin/env bash
# Downloads updates to pacman cache daily for quicker Syus
checkupdates -d
